# metal/monolithic/application/overlay/etc/supervisor.d/50_sound.ini
[program:sound]
priority=50
command=/usr/bin/50_start_%(program_name)s.sh
# logging
stdout_logfile=%(ENV_XDG_RUNTIME_DIR)s/logs/log_%(program_name)s
stdout_logfile_maxbytes=50MB
stderr_logfile=%(ENV_XDG_RUNTIME_DIR)s/logs/log_%(program_name)s
stderr_logfile_maxbytes=50MB
# options
stopasgroup=true
stopsignal=INT
environment=
  XDG_RUNTIME_DIR="%(ENV_XDG_RUNTIME_DIR)s"
# start
autostart=true
autorestart=true
